# config/systemd/user/bluetooth-tray.service
[Unit]
Description=Bluetooth tray (Mint)

[Service]
ExecStart=/usr/bin/blueberry-tray
