# config/systemd/user/brave.service
[Unit]
Description=Best web browser in the world

[Service]
ExecStart=/usr/bin/brave
