# config/systemd/user/network-tray.service
[Unit]
Description=Network tray (nm-applet)

[Service]
ExecStart=/usr/bin/nm-applet
