# config/systemd/user/notification.service
[Unit]
Description=Notification (XFCE4)

[Service]
ExecStart=/usr/lib/xfce4/notifyd/xfce4-notifyd
