# config/systemd/user/numlockx.service
[Unit]
Description=Numlock at startup

[Service]
ExecStart=/usr/bin/numlockx
