# config/systemd/user/polkit.service
[Unit]
Description=Polkit Agent (KDE)

[Service]
ExecStart=/usr/lib/polkit-kde-authentication-agent-1
