# config/systemd/user/power-management.service
[Unit]
Description=Power Management (XFCE4)

[Service]
ExecStart=/usr/bin/xfce4-power-manager --daemon
