# config/systemd/user/slack.service
[Unit]
Description=Worst communication tool ever

[Service]
ExecStart=/usr/bin/slack
