# config/systemd/user/sound-tray.service
[Unit]
Description=Volume and audio device tray (pasystray)

[Service]
ExecStart=/usr/bin/pasystray
