# config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper (variety)

[Service]
ExecStart=/usr/lib/varietry
