# config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper (variety)

[Service]
ExecStart=feh --recursive --bg-fill --randomize /usr/share/wallpapers
